# com.dailybriefing.generate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dailybriefing.generate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd /Users/sizhuang/Desktop/daily-briefing/backend &amp;&amp; /Users/sizhuang/anaconda3/condabin/conda run -n daily-briefing python generate_briefing.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sizhuang/Desktop/daily-briefing/output/generate.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sizhuang/Desktop/daily-briefing/output/generate.err</string>

    <key>WorkingDirectory</key>
    <string>/Users/sizhuang/Desktop/daily-briefing/backend</string>
</dict>
</plist>
